# ops/gtkram-backup.service
[Unit]
Description=GT Kram Backup-Script

[Service]
Type=oneshot
ExecStart=/opt/backup/gtkram-backup.sh
